# backend/deploy/supervisor/miioo-backend.conf
[program:miioo-web]
directory=/www/wwwroot/miiooaib.com/backend
command=/bin/bash ./脚本/start_prod.sh
user=www
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/www/wwwroot/miiooaib.com/backend/logs/supervisor/web.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
redirect_stderr=true
; MEDIA-BPLUS 首轮主域复用默认口径：
; - 只有 www.miiooai.com 一个后端域名时，先用 /media/origin 与 /media/cdn 承接对象存储地址
; - 首次上云建议先保留 MEDIA_STORAGE_MODE=local，待 COS/CDN secret、bucket 与 Nginx 回源都填完后，再切到 hybrid / object_storage
environment=PROD_VENV_DIR="/www/wwwroot/miiooaib.com/backend/.venv",APP_ENV="production",BACKGROUND_JOB_EXECUTION_MODE="queue",WORKERS="2",HOST="0.0.0.0",PORT="8000",FORWARDED_ALLOW_IPS="*",MEDIA_STORAGE_MODE="local",MEDIA_PUBLIC_BASE_URL="https://www.miiooai.com/media/origin",MEDIA_CDN_BASE_URL="https://www.miiooai.com/media/cdn",MEDIA_OBJECT_STORAGE_PROVIDER="tencent_cos",MEDIA_OBJECT_STORAGE_REGION="",MEDIA_OBJECT_STORAGE_SECRET_ID="",MEDIA_OBJECT_STORAGE_SECRET_KEY="",MEDIA_OBJECT_STORAGE_BUCKET_RAW="",MEDIA_OBJECT_STORAGE_BUCKET_PREVIEW="",MEDIA_OBJECT_STORAGE_BUCKET_DERIVED="",MEDIA_OBJECT_STORAGE_BUCKET_HLS="",MEDIA_OBJECT_STORAGE_BUCKET_PRIVATE_DOWNLOAD="",MEDIA_OBJECT_STORAGE_KEY_PREFIX_RAW="raw",MEDIA_OBJECT_STORAGE_KEY_PREFIX_PREVIEW="preview",MEDIA_OBJECT_STORAGE_KEY_PREFIX_DERIVED="derived",MEDIA_OBJECT_STORAGE_KEY_PREFIX_HLS="hls",MEDIA_OBJECT_STORAGE_KEY_PREFIX_PRIVATE_DOWNLOAD="private-download"

[program:miioo-worker]
directory=/www/wwwroot/miiooaib.com/backend
command=/bin/bash ./start_worker.sh
user=www
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/www/wwwroot/miiooaib.com/backend/logs/supervisor/worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
redirect_stderr=true
environment=PROD_VENV_DIR="/www/wwwroot/miiooaib.com/backend/.venv",APP_ENV="production",BACKGROUND_JOB_EXECUTION_MODE="queue",MEDIA_STORAGE_MODE="local",MEDIA_PUBLIC_BASE_URL="https://www.miiooai.com/media/origin",MEDIA_CDN_BASE_URL="https://www.miiooai.com/media/cdn",MEDIA_OBJECT_STORAGE_PROVIDER="tencent_cos",MEDIA_OBJECT_STORAGE_REGION="",MEDIA_OBJECT_STORAGE_SECRET_ID="",MEDIA_OBJECT_STORAGE_SECRET_KEY="",MEDIA_OBJECT_STORAGE_BUCKET_RAW="",MEDIA_OBJECT_STORAGE_BUCKET_PREVIEW="",MEDIA_OBJECT_STORAGE_BUCKET_DERIVED="",MEDIA_OBJECT_STORAGE_BUCKET_HLS="",MEDIA_OBJECT_STORAGE_BUCKET_PRIVATE_DOWNLOAD="",MEDIA_OBJECT_STORAGE_KEY_PREFIX_RAW="raw",MEDIA_OBJECT_STORAGE_KEY_PREFIX_PREVIEW="preview",MEDIA_OBJECT_STORAGE_KEY_PREFIX_DERIVED="derived",MEDIA_OBJECT_STORAGE_KEY_PREFIX_HLS="hls",MEDIA_OBJECT_STORAGE_KEY_PREFIX_PRIVATE_DOWNLOAD="private-download"
